# docker-crontab-python/cronpy/crontab
* * * * * root sh /cronpy/cronjob.sh 1 >> /cronpy/cronjob.log 2>&1
0 */2 * * * root sh /cronpy/cronjob.sh 2 >> /cronpy/cronjob.log 2>&1
# an empty line is required by crintab, don't remove it.
